# systemd/fresh-momentum.service
[Unit]
Description=Huragan Core - Fresh Momentum Capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/huragan_core
Environment=FRESH_MOMENTUM_CAPTURE=only
Environment=PAPER_MODE=true
Environment=LIVE_ARMED=false
ExecStart=/opt/huragan_core/target/release/huragan_core
Restart=always
RestartSec=5

# hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
